# flatpak/profile.d/flatpak-variables.sh
#!/usr/bin/env bash

export $XDG_DATA_DIRS="/var/lib/flatpak/exports/share:/home/ivan/.local/share/flatpak/exports/share:${XDG_DATA_DIRS}"
