# systemd.config/user/password-store-sync.timer
[Unit]
Description=Sync pasword-store git repository every hour

[Timer]
OnCalendar=*:0/5

[Install]
WantedBy=timers.target
