# systemd.config/user/password-store-sync.service
[Unit]
Description=Sync password-store git repository

[Service]
Type=simple
ExecStart=/bin/bash -c 'pass git pull --rebase && pass git push'
